# system_files/desktop/shared/etc/systemd/user/obex.service.d/override.conf
[Service]
ExecStart=
ExecStart=/usr/libexec/bluetooth/obexd --symlinks
