# system_files/desktop/shared/usr/lib/systemd/user/container-update.timer
[Unit]
Description=Timer for updating userspace containers

[Timer]
OnStartupSec=1m
OnUnitInactiveSec=1d

[Install]
WantedBy=timers.target
